# systemd/user/cronus-daily.service
[Unit]
Description=Cronus-Daily
[Service]
Type=oneshot
ExecStart=/home/robert/scripts/cronus-daily.sh
StandardOutput=journal
StandardError=journal
