# systemd/user/cronus-daily.timer
[Unit]
Description=Cronus PV-Forschung - Täglich um 04:00 Uhr
After=network.target

[Timer]
OnBootSec=5min
Unit=cronus-daily.service
Persistent=true

[Install]
WantedBy=timers.target
